# files/kolla-veth-XXX.service
[Unit]
After=network.service

[Service]
ExecStart=/usr/local/bin/veth-__NETNAME__.sh

[Install]
WantedBy=default.target
